# scripts/com.aurelia.agent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aurelia.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/.aurelia</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/.aurelia/logs/aurelia.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.aurelia/logs/aurelia.stderr.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
